# ci/systemd/group02-deploy.timer
[Unit]
Description=Run the delivery deploy hook every ~2 minutes

[Timer]
# OnBootSec: fire the first run 1 min after the timer starts (login for a --user timer),
# giving the machine a moment to settle. OnUnitActiveSec: then re-fire every 2 min,
# measured from when the *service last became active* — so the interval is between runs,
# not overlapping with a run still in progress. AccuracySec loosens the scheduler's
# precision to save wakeups; 15s is plenty tight for a 2-min poll.
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=15s
Unit=delivery-deploy.service

[Install]
# Hooking into timers.target is what makes `systemctl enable` start this timer
# automatically — without an [Install] section the timer couldn't be enabled.
WantedBy=timers.target
