# ci/systemd/group02-deploy.service
[Unit]
Description=delivery delivery-risk deploy hook (pull main + deploy on change)
After=network-online.target

[Service]
# oneshot = run the command to completion and exit (not a long-lived daemon). This is a
# "do one deploy check, then stop" job; the paired .timer is what makes it recur. Pairing
# a oneshot service with a timer is the systemd idiom for a cron-style periodic task.
Type=oneshot
# %h = the user's home (this is a --user unit). The hook is self-contained and reads its
# own config, so the unit stays free of any secrets or per-environment paths.
ExecStart=/usr/bin/env bash %h/project/ci/deploy_hook.sh
# Kill a hung run after 10 min so a stuck test/model-load can't block future ticks
# forever; the hook's flock separately prevents two runs from overlapping.
TimeoutStartSec=600
